# crontab
*/10 * * * * /config/run.sh >> /config/log/cron.log 2>&1
